# scripts/com.epicinvest.pairwise-shadow-trader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.epicinvest.pairwise-shadow-trader</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jsa/work/epic_invest/scripts/pairwise_shadow_launchd_entry.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jsa/work/epic_invest</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/tmp/epic-invest-pairwise-shadow.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/epic-invest-pairwise-shadow.log</string>
</dict>
</plist>
